# modules/ubuntu/files/apt-synch.cron
@midnight root /usr/local/sbin/apt-synch
